# supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=java -jar /opt/leshan/leshan-server-demo.jar

[program:client]
command=java -jar /opt/leshan/leshan-client-demo.jar
